# sparse/usr/lib/systemd/system/ofono-pre.service
[Service]
TimeoutStartSec=infinity
ConditionPathExists=/dev/socket/qmux_radio/qcril_radio_config0
ExecStart=/usr/bin/sleep 1000
RemainAfterExit=yes

